# Docker/8/supervisord/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=nobody:nogroup

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:php]
command=php -d variables_order=EGPCS /var/www/html/artisan serve --host=0.0.0.0 --port=80
user=www-data

[program:queue]
command=php -d variables_order=EGPCS /var/www/html/artisan queue:work --sleep=3 --timeout=10 --tries=2 --delay=600
user=www-data


[include]
files=/etc/supervisord.d/*.conf
